# .config/cron/crontab
# Update RSS
0 */5 * * * ~/.config/cron/rssup.sh blogs soft yt 2>&1 >/dev/null
