# .config/cron/crontab
# Update RSS
0 */2 * * * ~/bin/rssup blogs soft videos 2>&1 >/dev/null; pkill -HUP status
0 7 * * 1 ~/bin/sfeed_clean
