# .config/cron/crontab
# Update RSS
0 */2 * * * ~/bin/rssup blogs soft videos 2>&1 >/dev/null; pkill -RTMIN+3 status
